# mnt-share.automount
[Unit]
Description=Automount sharename
Requires=network-online.target
After=network-online.service

[Automount]
Where=/mnt/sharename
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target